# ansible/roles/install/templates/tendermint.service.jinja2
[Unit]
Description={{ service }}
Requires=network-online.target
After=network-online.target

[Service]
Restart=on-failure
User={{ user }}
Group={{ service }}